# rust-actix-watcher.path
[Path]
PathModified=/home/rust-actix/rust-actix

[Install]
WantedBy=multi-user.target